# systemd/battery-logger.service
[Unit]
Description=Battery Logger

[Service]
Type=simple
ExecStart=/home/lukasz/.local/bin/BatteryLogger

[Install]
RequiredBy=battery-logger.timer